# .config/pacpac/todo/rootfiles/etc/systemd/system/daily-reboot.timer
[Unit]
Description=daily reboot for remote reasons

[Timer]
OnCalendar=*-*-* 04:19
Unit=reboot.target

[Install]
WantedBy=timers.target
