# .config/pacpac/todo/rootfiles/etc/systemd/system/wondershaper.service.d/10-no-fail-on-exit.conf
[Service]
ExecStop=
ExecStop=-/usr/sbin/wondershaper -c -a $IFACE
